# playbook/templates/zfs_exporter.service.j2
[Unit]
Description=zfs_exporter 

[Service]
Type=simple
ExecStart=/usr/local/bin/zfs_exporter
Restart=always

[Install]
WantedBy=multi-user.target
